# runit/ws_rtsp/run
#!/bin/sh -e
exec 2>&1
exec chpst -u root websock_rtsp_proxy --config /etc/ws_rtsp.ini

